# backup.service
[Unit]
Description=Backup

[Service]
ExecStart=/home/ch/src/backup/backup.sh
StandardOutput=syslog
StandardError=syslog
