# backup.timer
[Unit]
Description=Backups

[Timer]
OnBootSec=15min
OnUnitInactiveSec=1d

[Install]
WantedBy=timers.target
