# backup.timer
[Unit]
Description=Backups

[Timer]
OnBootSec=15min
OnUnitInactiveSec=2h

[Install]
WantedBy=timers.target
