# offline-backup.timer
[Unit]
Description=Offline Backups

[Timer]
OnCalendar=*-*-* 02:17:00
Persistent=true

[Install]
WantedBy=timers.target
